# etc/profile.d/cvs.csh
# change default from rsh to ssh for cvs command
setenv CVS_RSH "${CVS_RSH-ssh}"