# etc/profile.d/libaccounts-glib.sh
# Set libaccounts-glib environment variables
export AG_SERVICES=/usr/share/accounts/services
export AG_SERVICE_TYPES=/usr/share/accounts/service-types
export AG_PROVIDERS=/usr/share/accounts/providers
