# etc/profile.d/vim.csh
[ -x //usr/bin/id ] || exit
[ `//usr/bin/id -u` -gt 100 ] && alias vi vim
