# etc/udev/rules.d/20-sdcard.rules
KERNEL=="mmcblk0p1", ACTION=="add",  RUN+="/bin/mount /dev/mmcblk0p1"
KERNEL=="mmcblk0p1", ACTION=="remove",  RUN+="/bin/umount /dev/mmcblk0p1"
