# usr/share/gvfs/mounts/burn.mount
[Mount]
Type=burn
Exec=/usr/libexec/gvfsd-burn
AutoMount=true

