# usr/share/gvfs/mounts/computer.mount
[Mount]
Type=computer
Exec=/usr/libexec/gvfsd-computer
AutoMount=true

