# usr/share/gvfs/mounts/dav+sd.mount
[Mount]
Type=dav+sd;davs+sd
Exec=/usr/libexec/gvfsd-dav
AutoMount=false
