# usr/share/gvfs/mounts/dav.mount
[Mount]
Type=dav;davs
Exec=/usr/libexec/gvfsd-dav
AutoMount=false
