# usr/share/gvfs/mounts/dns-sd.mount
[Mount]
Type=dns-sd
Exec=/usr/libexec/gvfsd-dnssd
AutoMount=true
