# usr/share/gvfs/mounts/ftp.mount
[Mount]
Type=ftp
Exec=/usr/libexec/gvfsd-ftp
AutoMount=false
Scheme=ftp
DefaultPort=21
HostnameIsInetAddress=true
