# usr/share/gvfs/mounts/http.mount
[Mount]
Type=http
Exec=/usr/libexec/gvfsd-http
AutoMount=true
DBusName=org.gtk.vfs.mountpoint.http
