# usr/share/gvfs/mounts/localtest.mount
[Mount]
Type=localtest
Exec=/usr/libexec/gvfsd-localtest
AutoMount=false
