# usr/share/gvfs/mounts/network.mount
[Mount]
Type=network
Exec=/usr/libexec/gvfsd-network
AutoMount=true
