# etc/xdg/autostart/DHAVN_VideoThumbNailer.desktop
[Desktop Entry]
Exec=/app/bin/DHAVN_VideoThumbNailer
