# etc/xdg/autostart/DHAVN_TrackerExtractorPlugin.desktop
[Desktop Entry]
Exec=/app/bin/DHAVN_TrackerExtractorPlugin
