# etc/xdg/autostart/messageserver.desktop
[Desktop Entry]
Exec=/usr/bin/messageserver
X-Meego-Priority=Low
OnlyShowIn=X-MEEGO-HS;
