# etc/xdg/autostart/msyncd.desktop
[Desktop Entry]
Exec=/usr/bin/msyncd
Name=Synchronization daemon
Type=Application
